# resume-service/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi --ini uwsgi.ini
stopsignal=QUIT

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
