# resume-service/deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi --ini /etc/uwsgi
stopsignal=QUIT

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
